# local/cron/crontab
# crontab
# see /var/log/cron.log
*/5 * * * * $HOME/dotfiles/local/cron/gms-daemon.sh
